# config/systemd/user/server-backup.timer
[Unit]
Description=Scheduled backup of minecraft server

[Timer]
OnCalendar=*-*-* 23:30
Persistent=false

[Install]
WantedBy=timers.target
